# cron/crontab/duclos
*/25 * * * * /home/duclos/.scripts/cron/cron_cleanup.sh
