# cron/crontab/root
5 * */5 * * /usr/bin/pacman -Syuw --noconfirm
